# s6-rc.d/sshd/run
#!/command/execlineb -P
exec 2>&1
s6-setsid -q -- /bin/sh -c "echo Docker! | su -c '/usr/sbin/sshd -D -e'"
