# s6-rc.d/sshd/run
#!/command/execlineb -P
s6-sudo -q -- /usr/sbin/sshd -D -e
