# s6-rc.d/dropbear-log/run
#!/command/execlineb -P
exec logutil-service /var/log/dropbear
