# s6-rc.d/fail2ban/run
#!/command/execlineb -S1
/usr/bin/fail2ban-client -x -f start
