# s6-rc.d/fail2ban/run
#!/command/execlineb -P
s6-sudo -q -- fail2ban-server -f
