# s6-rc.d/fail2ban/run
#!/command/execlineb -S1
rm -f /var/run/fail2ban/fail2ban.sock
exec /usr/bin/fail2ban-server -f
