# s6-rc.d/dropbear/run
#!/command/execlineb -P
/usr/sbin/dropbear -F -E -p 22 -P /var/run/dropbear.pid