# s6-rc.d/dropbear/run
#!/command/execlineb -P
fdmove -c 2 1
/usr/sbin/dropbear -s -F -E -p 10022 -P /var/run/dropbear.pid
